# volumio/lib/systemd/system/wac.service
[Unit]
Description=WAC
After=avahi-daemon.service
Wants=avahi-daemon.service
BindsTo=sys-subsystem-net-devices-wlan0.device
After=sys-subsystem-net-devices-wlan0.device
[Service]
Type=simple
ExecStartPre=/bin/hotspot.sh
ExecStart=/usr/local/bin/WACServer
RestartSec=2
StandardOutput=syslog
StandardError=syslog
[Install]
WantedBy=multi-user.target

